# scripts/systemd/cc-live-dashboard.service
[Unit]
Description=Neuract OS Live Dashboard (OPC UA + Modbus Monitor)
After=cc-live-simulator.service network.target
Wants=cc-live-simulator.service

[Service]
Type=simple
User=rohith
WorkingDirectory=/home/rohith/desktop/CommandCenter/rl_training_data/dashboard
ExecStart=/home/rohith/.pyenv/versions/3.11.9/bin/python3 -m uvicorn server:app --host 0.0.0.0 --port 8055
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
